# app/crontab
00 00 * * * /home/pa4373/projects/midnight-dvr/app/record.py
